# templates/default/solr.upstart.conf.erb
#!upstart
description "Solr Server"

env SOLR_DIR=<%= node[:solr][:dir] %>

start on startup
stop on shutdown

respawn

exec sudo sh -c "cd $SOLR_DIR/example && java -jar start.jar"